# infrastructure/systemd/veridian-api.service
[Unit]
Description=Veridian API (FastAPI)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/Veridian/apps/api
EnvironmentFile=/opt/Veridian/.env
Environment=PYTHONPATH=/opt/Veridian/apps/api/src
ExecStart=/opt/Veridian/apps/api/.venv/bin/uvicorn veridian_api.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
